# .config/systemd/user/gpg-agent.service.d/override.conf
[Service]
Environment=GNUPGHOME=/home/rx14/.config/gnupg
