# .config/systemd/user/gpg-agent.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gpg-agent --supervised --options /home/rx14/.config/gnupg/gpg-agent.conf
Environment=GNUPGHOME=/home/rx14/.config/gnupg
