# .config/systemd/user/gpg-agent.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=%t/gnupg/d.k6gwzj1wdey13xzbrkhouced/S.gpg-agent
